# origin/crontab
# Attempt to renew letsencrypt certs
0 8 * * 1 /root/cert_renew
